# log-viewer-cli/bin/etc/linux-upstart/system/log-viewer.conf
description "LogViewer is a Web application for monitoring server logs in real-time in the browser."

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

respawn

setuid "$USER"
setgid "$USER"

chdir /path/to/log-viewer/

exec /bin/java -ea -Dlogback.configurationFile=./logback.xml -Dlog-viewer.config-file=config.conf -jar lib/log-viewer-cli-${project.version}.jar startup